# deploy/seele-backend.service
[Unit]
Description=Seele Backend Service
After=network.target

[Service]
Type=simple
User=<SERVICE_USER>
WorkingDirectory=<DEPLOY_DIR>/seele-backend
Environment=DEPLOY_ENV=prod
EnvironmentFile=<DEPLOY_DIR>/seele-backend/.env
ExecStart=<DEPLOY_DIR>/seele-backend/.venv/bin/python -m uvicorn app.main:app --host 0.0.0.0 --port 9000 --workers 1
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
